# templates/druple.crontab
 0 * * * * sudo -u www-data sh -c "cd WEBROOT/live  && ./console.php cron:run" >> LOGDIR/cron/live.log  2>&1
15 * * * * sudo -u www-data sh -c "cd WEBROOT/stage && ./console.php cron:run" >> LOGDIR/cron/stage.log 2>&1
30 * * * * sudo -u www-data sh -c "cd WEBROOT/dev   && ./console.php cron:run" >> LOGDIR/cron/dev.log   2>&1
